# etc/systemd/system/update-esp.path
[Unit]
Description=Copy kernel and initramfs image to EFI system partition

[Path]
PathChanged=/boot/initramfs-linux-fallback.img

[Install]
WantedBy=multi-user.target
WantedBy=system-update.target 
